# xbee-monitor.upstart.conf
# XBee monitor service

start on runlevel [2345]
stop on runlevel [!2345]

script
    exec /usr/sbin/xbee-monitor
end script
